# ops/systemd/alarisa@.service
[Unit]
Description=Alarisa personal assistant for Linux user %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/alarisa/app
Environment=HOME=/home/%i
Environment=PATH=/home/%i/.local/bin:/home/%i/.local/npm/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillMode=mixed
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
UMask=0077

# The application loads HOST, PORT, and SERVER_TYPE from the project-root .env.
# Keep HOST=127.0.0.1 for the Apache reverse-proxy deployment and verify it
# with ss after every service configuration change.

[Install]
WantedBy=multi-user.target
